# .config/autostart/fusuma.desktop
 [Desktop Entry]
Name=fusuma
Comment=Fusuma Multitouch Gestores on Linux
Exec=/usr/bin/fusuma
Terminal=false
Type=Application
